# etc/client-supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[program:inetd]
command=/usr/sbin/inetd -i /etc/client-inetd.conf
autorestart=true
autostart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:syslogd]
command=/usr/sbin/syslogd -n
autorestart=true
autostart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
